# app/templates/crontab
0 3 * * * root python /app/scripts/scheduler.py >> /var/log/cron.log 2>&1